# aurora-app.service
[Unit]
Description=Aurora Hardware Dashboard (Electron App)
After=aurora-metrics.service
Wants=aurora-metrics.service

[Service]
Type=simple
User=pcrm
WorkingDirectory=/projects/aurora-v1
Environment=DISPLAY=:0.1
Environment=HYTE_OUTPUT=DP-0
ExecStart=/usr/local/bin/npm run electron
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
